# tym-daemon.service.in
[Unit]
Description=tym daemon

[Service]
Type=simple
ExecStart=@prefix@/bin/tym --daemon

[Install]
WantedBy=graphical.target
