# tym-daemon.service
[Unit]
Description=tym daemon

[Service]
Type=simple
ExecStart=/usr/bin/tym --daemon

[Install]
WantedBy=default.target
